# config/crontab
# OSM planet files are updated weekly
0 21 * * 4 /usr/bin/unbuffer /usr/bin/python -u /home/ubuntu/osm-places/scripts/fetcher.py >> /home/ubuntu/log/fetcher.log
