# config/crontab
# OSM planet files are updated weekly
0 21 * * 4 /usr/bin/python -u /usr/local/osm-places/scripts/fetcher.py >> /var/log/fetcher.log
